# periodic/daily/disable-wifi
#!/usr/bin/env sh

curl -s --header "Content-Type: application/json" --request POST --data "0" $API_URL
